# scripts/systemd/gaia-cluster.service
[Unit]
Description=GAIA Cluster — production stack + HA failover + doctor verification
After=gaia-warm-pool-seed.service docker.service network-online.target
Requires=gaia-warm-pool-seed.service docker.service
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
User=azrael
Group=azrael
Environment=TERM=xterm-256color
Environment=PATH=/home/azrael/.nvm/versions/node/v25.4.0/bin:/usr/local/bin:/usr/bin:/bin
WorkingDirectory=/gaia/GAIA_Project

ExecStart=/gaia/GAIA_Project/scripts/gaia-startup.sh

# Stop HA candidates first, then production stack
ExecStop=/gaia/GAIA_Project/scripts/ha_stop.sh
ExecStop=/usr/bin/docker compose -f /gaia/GAIA_Project/docker-compose.yml down -t 20

TimeoutStartSec=900
TimeoutStopSec=90

[Install]
WantedBy=multi-user.target
